# supervisor_course15.conf
[program:course15]
command = /bin/bash -c 'cd var/www/fors/course15 && source bin/activate && python3 app.py' 
stopsignal=INT
stopasgroup=true
killasgroup=true

